# supervisord/cdp-proxy.conf
[program:cdp-proxy]
environment=HOME="/home/%(ENV_USER)s",USER="%(ENV_USER)s"
command=/usr/bin/socat tcp-listen:9224,fork,reuseaddr tcp:127.0.0.1:9223
autorestart=true
priority=850
user=%(ENV_USER)s
stdout_logfile=/var/log/neko/cdp-proxy.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2
redirect_stderr=true
